# deploy/personal-map-track.service
[Unit]
Description=Personal Map Track Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/src/gerenwangzhan-main
Environment=NODE_ENV=production
Environment=TRACK_PORT=3001
Environment=TRACK_HOST=127.0.0.1
Environment=TRACK_LOG_PATH=/var/log/personal_site/track.log
ExecStart=/usr/bin/node /opt/src/gerenwangzhan-main/server/track/track-server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

